# luci-app-jsholat/root/etc/init.d/jsholat-bot
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=01
DEPENDS="network"

# Konfigurasi
PID_FILE="/var/run/jsholat-bot.pid"
BOT_CHECKER_SCRIPT="/usr/share/jsholat/bot-checker.sh"
BOT_CHECKER_INTERVAL="*/2 * * * *"

setup_bot_checker_cron() {
    if ! crontab -l | grep -q "$BOT_CHECKER_SCRIPT"; then
        (crontab -l 2>/dev/null; echo "$BOT_CHECKER_INTERVAL $BOT_CHECKER_SCRIPT") | crontab -
        echo "✅ Cron bot checker ditambahkan"
    fi
}

remove_bot_checker_cron() {
    crontab -l | grep -v "$BOT_CHECKER_SCRIPT" | crontab -
    echo "✂️ Cron bot checker dihapus"
}

boot() {
    start_service
}

start_service() {
    echo "⚙️ Memulai jsholat-bot..."
    
    mkdir -p /etc/jsholat
    chmod 755 /etc/jsholat
    
    setup_bot_checker_cron

    procd_open_instance
    procd_set_param command /usr/bin/jsholat-bot /etc/jsholat/last_update_id
    procd_set_param respawn 300 5 10
    procd_set_param env HOME=/tmp
    procd_set_param limits nofile="unlimited"
    procd_set_param pidfile "$PID_FILE"
    procd_close_instance
    
    echo "✅ jsholat-bot started"
}

stop_service() {
    echo "🔌 Menghentikan jsholat-bot..."
    
    # Hanya matikan via procd (ini aman)
    procd_kill jsholat-bot 2>/dev/null
    
    # Hapus file PID dan lock
    rm -f "$PID_FILE" "/var/run/jsholat-bot.lock" 2>/dev/null
    
    # JANGAN pakai pgrep/kill karena bisa mematikan proses restart itu sendiri!
    # Biarkan procd yang handle
    
    remove_bot_checker_cron
    
    echo "✅ jsholat-bot stopped"
}

# PERBAIKAN: Fungsi restart yang benar
restart() {
    echo "♻️ Restarting jsholat-bot..."
    
    # Stop service via procd (tanpa membunuh diri sendiri)
    procd_kill jsholat-bot 2>/dev/null
    
    # Tunggu sebentar
    sleep 2
    
    # Start service
    start_service
    
    echo "✅ jsholat-bot restarted"
}

service_triggers() {
    procd_add_reload_trigger "jsholat"
}

# Fungsi status sederhana
status() {
    if procd_running jsholat-bot; then
        echo "running"
        return 0
    else
        echo "inactive"
        return 1
    fi
}